# exalt/data/daemon/init.d/gentoo
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later


depend() {
    need dbus hald
}


start() {
    ebegin "Starting Exalt daemon"
        start-stop-daemon --start --exec /usr/local/sbin/exalt-daemon \
        --pidfile /var/run/exaltd.pid
        eend $?
}

stop() {
    ebegin "Stopping Exalt daemon"
        start-stop-daemon --stop --exec /usr/local/sbin/exalt-daemon  \
        --pidfile /var/run/exaltd.pid
        eend $?
}

restart() {
    svc_stop
        svc_start
}

